# supervisor/show_env.conf
[program:show_env]
directory=/opt/example-python/src
command=python show_env.py 80
